# stats_api/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn -w 4 -b 0.0.0.0:5000 wsgi:app
directory=/app
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
